# profile.d/misc.sh
export EDITOR=sandy
export PAGER=less
export LESS="-RS"
