# profile.d/misc.sh
export PATH="$HOME/.local/bin:$PATH"
export EDITOR=sandy
export PAGER=less
export LESS="-RS"
