# profile.d/rust.sh
export RUSTC_BOOTSTRAP=1
export RUSTC_WRAPPER=sccache
export RUSTFLAGS="-Abad-style -Coverflow-checks=no"
