# profile.d/rust.sh
export RUSTC_BOOTSTRAP=1
export RUSTFLAGS="-Abad-style"
